# dev-android/adbd/files/adbd-init
#!/sbin/runscript

PIDFILE=/var/run/adbd.pid

depend() {
	use logger
}

start() {
	ebegin "Starting adbd"
		echo -n >/sys/class/android_usb/android0/enable "0" && \
		echo -n >/sys/class/android_usb/android0/idVendor "0B05" && \
		echo -n >/sys/class/android_usb/android0/idProduct "4C91" && \
		echo -n >/sys/class/android_usb/android0/functions "adb" && \
		echo -n >/sys/class/android_usb/android0/iManufacturer "Asus" && \
		echo -n >/sys/class/android_usb/android0/iProduct "Transformer Infinity" && \
		echo -n >/sys/class/android_usb/android0/iSerial "0123456789ABCDEF" && \
		echo -n >/sys/class/android_usb/android0/enable "1" && \
		start-stop-daemon --start --quiet --background \
			--make-pidfile --pidfile ${PIDFILE} \
			--exec /usr/bin/adbd
	eend $?
}

stop() {
	ebegin "Stopping adbd"
		start-stop-daemon --stop --quiet \
			--pidfile ${PIDFILE}
	eend $?
	echo -n >/sys/class/android_usb/android0/enable "0"
}
